# kalshi_weather_model/ops/launchd/com.kalshi-weather.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  TEMPLATE / REFERENCE ONLY.
  The working LaunchAgent is generated by ops/launchd/install_launchd.sh,
  which substitutes the real project directory and user paths at install
  time. Replace /path/to/KalshiBot below if using this file directly.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kalshi-weather.scheduler</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/KalshiBot/kalshi_weather_model/ops/run_scheduler.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/KalshiBot/kalshi_weather_model</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/path/to/KalshiBot/kalshi_weather_model/data/logs/launchd_scheduler.out.log</string>

  <key>StandardErrorPath</key>
  <string>/path/to/KalshiBot/kalshi_weather_model/data/logs/launchd_scheduler.err.log</string>
</dict>
</plist>
